# s6/debian-root/etc/services.d/unbound/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping Unbound"
killall -15 unbound
/usr/lib/unbound/package-helper chroot_teardown
